# systemd/autoupdate.service
[Unit]
Description = Arch Auto Update Service

[Service]
Type=oneshot
ExecStart=/home/pogoplug/autoupdate/autoupdate.sh
